# deploy/systemd/magpie-ingest@fred.timer
[Unit]
Description=Magpie FRED ingest — daily after US market close

[Timer]
OnCalendar=Mon..Fri *-*-* 19:00:00
Unit=magpie-ingest@fred.service
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
