# deploy/systemd/quantfoundry-rebuild-catalog.service
[Unit]
Description=Magpie signal catalog rebuild + MinIO mirror
Documentation=https://github.com/your-org/Magpie/blob/main/server/signals/rebuild-catalog.ts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/quantfoundry/Magpie
EnvironmentFile=/srv/quantfoundry/.env
# CATALOG_S3_PUSH=1 enables the post-rebuild aws s3 cp of portfolio.duckdb
# to s3://quantfoundry-data/duckdb/. Laptops fetch via scripts/fetch-catalog.sh.
Environment=CATALOG_S3_PUSH=1
Environment=CATALOG_DB_PATH=/srv/quantfoundry/portfolio.duckdb
ExecStart=/usr/bin/npx tsx server/signals/rebuild-catalog.ts /srv/quantfoundry/Magpie/data
StandardOutput=append:/srv/quantfoundry/logs/rebuild-catalog.log
StandardError=inherit
Nice=10

[Install]
WantedBy=default.target
