# deploy/systemd/quantfoundry-ingest@.service
[Unit]
Description=Magpie data ingest — %I
Documentation=https://github.com/your-org/Magpie/blob/main/data/CRON.md
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/srv/quantfoundry/Magpie
EnvironmentFile=/srv/quantfoundry/.env
ExecStart=/usr/bin/npm run --silent ingest -- --source %i
StandardOutput=append:/srv/quantfoundry/logs/ingest-%i.log
StandardError=inherit
Nice=10

[Install]
WantedBy=default.target
